# etc/service/npcd/run
#!/bin/bash -e

exec /usr/local/pnp4nagios/bin/npcd -f /usr/local/pnp4nagios/etc/npcd.cfg 2>&1
